# crontab
0 * * * * /icebreaker/scripts/refresh-token.sh